# orchestration/fanout/launchd/com.user.fanout-ccb-sync.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- re-adapt after a ccb update — WatchPaths watches ccb's VERSION file, a ccb update rewrite triggers a full adapt.
     install: replace the placeholders with your absolute paths, copy to ~/Library/LaunchAgents/, launchctl load -w it.
       __HOME__       → your $HOME (e.g. /Users/leo)
       __FANOUT__     → fanout driver absolute path (e.g. __HOME__/.claude/skills/fanout/fanout)
       __CCB_WORK__   → main ccb project root (e.g. __HOME__/Projects/ccb-test)
       __CCB_CLAUDE__ → cc-claude pool project root (e.g. __HOME__/Projects/ccb-claude-only)
       __CCB_VERSION_FILE__ → ccb runtime VERSION (the VERSION under the Install path from ccb version)
     adapt --apply kills each project's ccbd(clears old code, next ccb start of the daemon loads the new version). Trigger timing = you just ran `ccb update`,
     usually not mid fan-out, so kill is safe. For more safety swap `adapt --apply` in the command back to `check`(report only, does not touch ccbd). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.user.fanout-ccb-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>echo "===== $(date '+%Y-%m-%d %H:%M') ccb update→adapt ====="; CCB_WORK=__CCB_WORK__ CCB_CLAUDE=__CCB_CLAUDE__ __FANOUT__ ccb-sync adapt --apply</string>
  </array>

  <!-- ccb update rewrites VERSION → triggers -->
  <key>WatchPaths</key>
  <array>
    <string>__CCB_VERSION_FILE__</string>
  </array>

  <key>RunAtLoad</key><false/>
  <key>ThrottleInterval</key><integer>30</integer>

  <key>StandardOutPath</key><string>__HOME__/Library/Logs/fanout-ccb-sync.log</string>
  <key>StandardErrorPath</key><string>__HOME__/Library/Logs/fanout-ccb-sync.log</string>
</dict>
</plist>
